# pkg/rpm/glideinmonitor-indexer.cron
23 * * * * root [ -f /var/lock/subsys/glideinmonitor-indexer ] && sleep $((RANDOM \% 600)) && /usr/sbin/glideinmonitor-indexer
